# ops/launchd/net.daytradingbot.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>net.daytradingbot.weekly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/richardducat/GITHUB/daytradingbot.net/scripts/generate-weekly-recap.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>/Users/richardducat</string>
    <key>PATH</key><string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key><false/>
  <key>WorkingDirectory</key><string>/Users/richardducat/GITHUB/daytradingbot.net</string>
  <key>StartCalendarInterval</key>
  <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>/Users/richardducat/.openclaw/logs/daytradingbot-weekly.log</string>
  <key>StandardErrorPath</key><string>/Users/richardducat/.openclaw/logs/daytradingbot-weekly.log</string>
</dict>
</plist>
